# univention-openvpn/debian/univention-openvpn-server.postrm
#!/bin/sh

#DEBHELPER#

ui=06univention-openvpn-server-uninstall.uinst

if [ "$1" = "remove" ]; then
	. /usr/share/univention-lib/base.sh
	call_unjoinscript "$ui"
	invoke-rc.d univention-directory-listener restart
fi

exit 0
